# packages/ws-storage/ws_storage-0.6b20-py36-none-any.whl/ws_storage/ws_storage.service
[Unit]
Description=python binary storage service

[Service]
ExecStart=/usr/local/bin/ws_storage runserver simple --conf_dir /etc/ws_storage/conf

[Install]
WantedBy=multi-user.target

