# packages/ws-sheets-server/ws_sheets_server-0.4.1b37.tar.gz/ws_sheets_server-0.4.1b37/ws_sheets_server/ws_sheets_server.service
[Unit]
Description=web sheets calculation server

[Service]
ExecStart=/usr/local/bin/ws_sheets_server runserver simple --conf_dir=/etc/ws_sheets_server/conf

[Install]
WantedBy=multi-user.target

