# packages/noolite-mqtt-webserver/noolite-mqtt-webserver-0.1.2.tar.gz/noolite-mqtt-webserver-0.1.2/noolite_mqtt_webserver/systemd_script/noolite_web_server.service
[Unit]
Description=NooLite WEB mqtt server for HomeKit
After=mosquitto.service

[Service]
Type=idle
ExecStart=/usr/local/bin/noolite_web_server

[Install]
WantedBy=multi-user.target
